# deploy/supervisor/laravel-worker-staging.conf
; =============================================================================
; Laravel Queue Worker - STAGING (Cloudways)
; =============================================================================
; Queue driver: redis (see QUEUE_CONNECTION in .env)
;
; Install on the Cloudways server at:
;   /etc/supervisor/conf.d/laravel-worker-staging.conf
;
; Before using this file, replace the placeholders below:
;   <STAGING_APP_FOLDER>  -> Cloudways application folder name
;                           (e.g. abcdefghij - found in Cloudways > Application
;                           > Access Details > "Application Folder")
;   <STAGING_APP_USER>    -> Cloudways application SSH/SFTP username
;                           (e.g. master_xxxxxxxx)
;   <PHP_BINARY>          -> Path to the PHP binary for this app
;                           (e.g. /usr/bin/php8.2). Confirm with:
;                             which php8.2
; =============================================================================

[program:laravel-worker-staging]
process_name=%(program_name)s_%(process_num)02d
command=php8.2 /home/master/applications/dguzgvtbzk/public_html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --backoff=5 --timeout=60 --queue=default
directory=/home/master/applications/dguzgvtbzk/public_html
autostart=true
autorestart=true
startsecs=1
startretries=3
stopasgroup=true
killasgroup=true
user=master_gjnpnxdvwb
numprocs=1
redirect_stderr=true
stdout_logfile=/home/master/applications/dguzgvtbzk/public_html/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
environment=LARAVEL_ENV="staging"
