# deploy/supervisor/laravel-worker-production.conf
; =============================================================================
; Laravel Queue Worker — PRODUCTION (Cloudways)
; =============================================================================
; Driver: Redis (set QUEUE_CONNECTION=redis in .env on the server)
;
; Install:
;   sudo cp laravel-worker-production.conf /etc/supervisor/conf.d/
;   (or paste contents into a new file there)
;
; Replace placeholders before activating:
;   <PROD_APP_FOLDER>  Application folder from Cloudways → Access Details
;   <PROD_APP_USER>    SSH/SFTP username for this application
;   <PHP_BINARY>       e.g. /usr/bin/php8.2  (run: which php8.2)
; =============================================================================

[program:laravel-worker-production]
process_name=%(program_name)s_%(process_num)02d
command=<PHP_BINARY> /home/master/applications/<PROD_APP_FOLDER>/public_html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --backoff=5 --timeout=60 --queue=default
directory=/home/master/applications/<PROD_APP_FOLDER>/public_html
autostart=true
autorestart=true
startsecs=1
startretries=3
stopasgroup=true
killasgroup=true
user=<PROD_APP_USER>
numprocs=2
redirect_stderr=true
stdout_logfile=/home/master/applications/<PROD_APP_FOLDER>/public_html/storage/logs/worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stopwaitsecs=3600
environment=LARAVEL_ENV="production"
